# dotboss@.service
[Unit]
Description=Watch dotfiles repository and git commit all changes using gitwatch.

[Service]
EnvironmentFile=%h/.config/dotboss/config
ExecStart=bash -c "/usr/local/bin/gitwatch -r ${DOT_REPO_REMOTE} -b ${DOT_REPO_BRANCH} -l ${USER} -s 60 ${DOT_REPO_DIR}"
ExecStop=/bin/true

[Install]
WantedBy=default.target
